# deploy/systemd/code-memory-reindex@.timer
[Unit]
Description=Daily code memory sync and reindex for %i

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
